# rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

#
# CSCS Note: This file is maintained via Dockerfile AND instruction
#


touch /var/lock/subsys/local

yum -y update ca-policy-egi-core

# END-OF-FILE
